# com.vaptanix.instagram-poster.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vaptanix.instagram-poster</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/harshakar/Documents/Instagram Automation/venv/bin/python</string>
        <string>/Users/harshakar/Documents/Instagram Automation/main.py</string>
        <string>--schedule</string>
        <string>--ai</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/harshakar/Documents/Instagram Automation/bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/harshakar/Documents/Instagram Automation/bot.error.log</string>
</dict>
</plist>
